# scripts/dm_event_systemd_red_hat.socket
[Unit]
Description=Device-mapper event daemon FIFOs
DefaultDependencies=no

[Socket]
ListenFIFO=/var/run/dmeventd-server
ListenFIFO=/var/run/dmeventd-client
SocketMode=0600

[Install]
WantedBy=sockets.target
