# scripts/lvm-devices-import.path
[Unit]
Description=lvm-devices-import to create system.devices
ConditionPathExists=!/etc/lvm/devices/system.devices

# /run/lvm/lvm-devices-import created by vgchange -aay <rootvg>

[Path]
PathExists=/run/lvm/lvm-devices-import
Unit=lvm-devices-import.service

[Install]
WantedBy=multi-user.target
