# scripts/lvm2_lvmetad_systemd_red_hat.service.in
[Unit]
Description=LVM2 metadata daemon
Requires=lvm2-lvmetad.socket
After=lvm2-lvmetad.socket
DefaultDependencies=no
Conflicts=shutdown.target

[Service]
Type=forking
NonBlocking=true
ExecStart=@sbindir@/lvmetad
ExecStartPost=@sbindir@/lvm pvscan --cache
ExecReload=@sbindir@/lvmetad -R
Environment=SD_ACTIVATION=1
Restart=on-abort
PIDFile=@DEFAULT_RUN_DIR@/lvmetad.pid

[Install]
WantedBy=sysinit.target
