# scripts/lvm2_cluster_activation_systemd_red_hat.service.in
[Unit]
Description=Clustered LVM volumes activation service
Requires=lvm2-clvmd.service
After=lvm2-clvmd.service lvm2-cmirrord.service
OnFailure=lvm2-clvmd.service
DefaultDependencies=no
Conflicts=shutdown.target

[Service]
Type=simple
RemainAfterExit=yes
EnvironmentFile=-@sysconfdir@/sysconfig/clvmd
ExecStart=@systemdutildir@/lvm2-cluster-activation activate
ExecStop=@systemdutildir@/lvm2-cluster-activation deactivate

[Install]
WantedBy=multi-user.target
