# scripts/lvm2_pvscan_systemd_red_hat@.service.in
[Unit]
Description=LVM2 PV scan on %I
Documentation=man:pvscan(8)
DefaultDependencies=no
After=lvm2-lvmetad.socket %i.device
Before=shutdown.target
Conflicts=shutdown.target

[Service]
Type=simple
ExecStart=@sbindir@/pvscan --cache --activate ay %I
