# scripts/lvm2_monitoring_systemd_red_hat.service.in
[Unit]
Description=Monitoring of LVM2 mirrors, snapshots etc. using dmeventd or progress polling
Requires=dm-event.socket
After=fedora-storage-init.service fedora-storage-init-late.service
Before=local-fs.target
DefaultDependencies=no
Conflicts=shutdown.target

[Service]
Type=oneshot
Environment=LVM_SUPPRESS_LOCKING_FAILURE_MESSAGES=1
ExecStart=@sbindir@/lvm vgchange --monitor y
ExecStop=@sbindir@/lvm vgchange --monitor n
RemainAfterExit=yes

[Install]
WantedBy=sysinit.target
Also=dm-event.socket
