# scripts/lvm2_clvmd_systemd_red_hat.service.in
[Unit]
Description=Clustered LVM daemon
Documentation=man:clvmd(8)
After=dlm.service corosync.service
Before=remote-fs.target
Requires=network.target dlm.service corosync.service
RefuseManualStart=true
RefuseManualStop=true
StopWhenUnneeded=true
DefaultDependencies=no
Conflicts=shutdown.target

[Service]
Type=forking
Environment=CLVMD_OPTS=-T30
EnvironmentFile=-@sysconfdir@/sysconfig/clvmd
ExecStart=@sbindir@/clvmd $CLVMD_OPTS
SuccessExitStatus=5
TimeoutStartSec=30
TimeoutStopSec=10
OOMScoreAdjust=-1000
Restart=on-abort
PIDFile=@CLVMD_PIDFILE@
