# scripts/lvm2_pvscan_systemd_red_hat@.service.in
[Unit]
Description=LVM2 PV scan on device %i
Documentation=man:pvscan(8)
DefaultDependencies=no
BindsTo=dev-block-%i.device
After=lvm2-lvmetad.socket
Before=shutdown.target
Conflicts=shutdown.target

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=@sbindir@/pvscan --cache --activate ay /dev/block/%i
ExecStop=@sbindir@/pvscan --cache %i
