# scripts/lvm2_lvmetad_systemd_red_hat.socket.in
[Unit]
Description=LVM2 metadata daemon socket
DefaultDependencies=no

[Socket]
ListenStream=@DEFAULT_RUN_DIR@/lvmetad.socket
SocketMode=0600

[Install]
WantedBy=sockets.target
